# PROJECT HEAD/sigma-file-manager-2/DeepSoul-Coder-1/scripts/deepseek-coder.service
[Unit]
Description=DeepSeek-Coder API Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=deepsoul
Group=deepsoul
WorkingDirectory=/opt/deepseek-coder
EnvironmentFile=-/opt/deepseek-coder/.env

# Environment variables
Environment=MODEL_SIZE=base
Environment=PORT=8000
Environment=LOG_LEVEL=INFO
Environment=MODEL_PATH=/opt/deepseek-coder/models
Environment=DATA_DIR=/opt/deepseek-coder/data
Environment=LOG_DIR=/opt/deepseek-coder/logs
Environment=TORCH_COMPILE=1
Environment=PYTORCH_CUDA_ALLOC_CONF=max_split_size_mb:128

# Uncomment to enable authentication
#Environment=ENABLE_AUTH=true
#Environment=API_KEYS=sk-your-api-key-1,sk-your-api-key-2

# Increase resource limits for large models
LimitNOFILE=65536
TimeoutStartSec=600
TimeoutStopSec=30

# Start the service
ExecStart=/usr/bin/python -m uvicorn app:app --host 0.0.0.0 --port ${PORT} --log-level ${LOG_LEVEL}

# Restart policy
Restart=on-failure
RestartSec=10s

# Security hardening
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
